# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8001
directory=/app/services/auth-service
environment=DATABASE_NAME="auth_db",PORT="8001"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:patient-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8002
directory=/app/services/patient-service
environment=DATABASE_NAME="patient_db",PORT="8002"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:doctor-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8003
directory=/app/services/doctor-service
environment=DATABASE_NAME="doctor_db",PORT="8003"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:appointment-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8004
directory=/app/services/appointment-service
environment=DATABASE_NAME="appointment_db",PORT="8004"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:telemedicine-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8005
directory=/app/services/telemedicine-service
environment=DATABASE_NAME="telemedicine_db",PORT="8005"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:payment-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8006
directory=/app/services/payment-service
environment=DATABASE_NAME="payment_db",PORT="8006"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:notification-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8007
directory=/app/services/notification-service
environment=DATABASE_NAME="notification_db",PORT="8007"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ai-symptom-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8008
directory=/app/services/ai-symptom-service
environment=DATABASE_NAME="ai_symptom_db",PORT="8008"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/services/gateway
environment=PORT="8000"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
