# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/services/api-gateway
autostart=true
autorestart=true
environment=PORT="8000",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth]
command=uvicorn app.main:app --host 127.0.0.1 --port 8001
directory=/app/services/auth-service
autostart=true
autorestart=true
environment=PORT="8001",DATABASE_NAME="auth_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:patient]
command=uvicorn app.main:app --host 127.0.0.1 --port 8002
directory=/app/services/patient-service
autostart=true
autorestart=true
environment=PORT="8002",DATABASE_NAME="patient_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:doctor]
command=uvicorn app.main:app --host 127.0.0.1 --port 8003
directory=/app/services/doctor-service
autostart=true
autorestart=true
environment=PORT="8003",DATABASE_NAME="doctor_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:appointment]
command=uvicorn app.main:app --host 127.0.0.1 --port 8004
directory=/app/services/appointment-service
autostart=true
autorestart=true
environment=PORT="8004",DATABASE_NAME="appointment_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:telemedicine]
command=uvicorn app.main:app --host 127.0.0.1 --port 8005
directory=/app/services/telemedicine-service
autostart=true
autorestart=true
environment=PORT="8005",DATABASE_NAME="telemedicine_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:payment]
command=uvicorn app.main:app --host 127.0.0.1 --port 8006
directory=/app/services/payment-service
autostart=true
autorestart=true
environment=PORT="8006",DATABASE_NAME="payment_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s",STRIPE_SECRET_KEY="%(ENV_STRIPE_SECRET_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:notification]
command=uvicorn app.main:app --host 127.0.0.1 --port 8007
directory=/app/services/notification-service
autostart=true
autorestart=true
environment=PORT="8007",DATABASE_NAME="notification_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ai-symptom]
command=uvicorn app.main:app --host 127.0.0.1 --port 8008
directory=/app/services/ai-symptom-service
autostart=true
autorestart=true
environment=PORT="8008",DATABASE_NAME="ai_symptom_db",MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
