# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:gateway]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/services/gateway
environment=PORT="8000",DATABASE_NAME="gateway_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:auth]
command=uvicorn app.main:app --host 127.0.0.1 --port 8001
directory=/app/services/auth-service
environment=PORT="8001",DATABASE_NAME="auth_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:patient]
command=uvicorn app.main:app --host 127.0.0.1 --port 8002
directory=/app/services/patient-service
environment=PORT="8002",DATABASE_NAME="patient_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:doctor]
command=uvicorn app.main:app --host 127.0.0.1 --port 8003
directory=/app/services/doctor-service
environment=PORT="8003",DATABASE_NAME="doctor_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:appointment]
command=uvicorn app.main:app --host 127.0.0.1 --port 8004
directory=/app/services/appointment-service
environment=PORT="8004",DATABASE_NAME="appointment_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:telemedicine]
command=uvicorn app.main:app --host 127.0.0.1 --port 8005
directory=/app/services/telemedicine-service
environment=PORT="8005",DATABASE_NAME="telemedicine_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:payment]
command=uvicorn app.main:app --host 127.0.0.1 --port 8006
directory=/app/services/payment-service
environment=PORT="8006",DATABASE_NAME="payment_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:notification]
command=uvicorn app.main:app --host 127.0.0.1 --port 8007
directory=/app/services/notification-service
environment=PORT="8007",DATABASE_NAME="notification_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:ai-symptom]
command=uvicorn app.main:app --host 127.0.0.1 --port 8008
directory=/app/services/ai-symptom-service
environment=PORT="8008",DATABASE_NAME="ai_symptom_db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
